# vs-code/rc.local
#!/bin/bash

su tklas -c '/usr/bin/code serve-web --host 127.0.0.1 --port 8888 --without-connection-token'